# ha-claude-addon/rootfs/etc/s6-overlay/s6-rc.d/claude-ttyd/finish
#!/usr/bin/with-contenv bashio

DEATH_COUNT_FILE="/data/watchdog/death-count"
WATCHDOG_ENABLED=$(bashio::config 'watchdog_enabled')

if [ "${WATCHDOG_ENABLED}" != "false" ]; then
    if [ "$1" -ne 0 ] 2>/dev/null; then
        CURRENT=$(cat "${DEATH_COUNT_FILE}" 2>/dev/null || echo "0")
        [[ "$CURRENT" =~ ^[0-9]+$ ]] || CURRENT=0
        echo $((CURRENT + 1)) > "${DEATH_COUNT_FILE}"
        bashio::log.warning "claude-ttyd crashed (exit=$1), death count: $((CURRENT + 1))"
    else
        rm -f "${DEATH_COUNT_FILE}"
        bashio::log.info "claude-ttyd exited cleanly, death count reset"
    fi
fi
