# ha-claude-addon/rootfs/etc/s6-overlay/s6-rc.d/claude-ttyd/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting claude-ttyd service..."

WATCHDOG_ENABLED=$(bashio::config 'watchdog_enabled')
DEATH_COUNT_FILE="/data/watchdog/death-count"

mkdir -p "/data/watchdog"

if [ "${WATCHDOG_ENABLED}" != "false" ]; then
    MAX_DEATH_TALLY=$(bashio::config 'max_death_tally')
    MAX_DEATH_TALLY=${MAX_DEATH_TALLY:-3}

    if [ ! -f "${DEATH_COUNT_FILE}" ]; then
        echo "0" > "${DEATH_COUNT_FILE}"
    fi

    CURRENT_DEATHS=$(cat "${DEATH_COUNT_FILE}" 2>/dev/null || echo "0")
    if ! [[ "$CURRENT_DEATHS" =~ ^[0-9]+$ ]]; then
        CURRENT_DEATHS=0
        echo "0" > "${DEATH_COUNT_FILE}"
    fi

    if [ "$CURRENT_DEATHS" -ge "$MAX_DEATH_TALLY" ]; then
        bashio::log.error "Service crashed ${CURRENT_DEATHS}/${MAX_DEATH_TALLY} times, stopping to prevent restart loop"
        rm -f "${DEATH_COUNT_FILE}"
        exit 1
    fi

    bashio::log.info "Watchdog: deaths=${CURRENT_DEATHS}/${MAX_DEATH_TALLY}"
fi

FONT_SIZE=$(bashio::config 'font_size')
THEME=$(bashio::config 'theme')
TMUX_ENABLED=$(bashio::config 'tmux_enabled')

TTYD_ARGS="-t fontSize=${FONT_SIZE:-14} -t theme=${THEME:-dark}"

bashio::log.info "Starting ttyd on port 7681 (tmux_enabled=${TMUX_ENABLED})"
exec ttyd ${TTYD_ARGS} -p 7681 -W /usr/bin/start-claude.sh
